# app/lightsail_ia/deploy/filmbot.service
[Unit]
Description=FilmBot Streamlit App
After=network.target

[Service]
Type=simple
User=filmbot
WorkingDirectory=/opt/filmbot/app/lightsail_ia
EnvironmentFile=/opt/filmbot/app/lightsail_ia/.env
ExecStart=/opt/filmbot/app/lightsail_ia/.venv/bin/streamlit run app.py \
  --server.port 8501 \
  --server.address 0.0.0.0 \
  --server.headless true \
  --browser.gatherUsageStats false \
  --server.maxUploadSize 5
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
